# config/qemu/prebuilt/rootfs/etc/inittab
::sysinit:/etc/rc.d/startup

# ::respawn:/bin/sh
::respawn:/sbin/getty -n -l /bin/sh -L 57600 ttyAMA0 vt100

::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc.d/shutdown
