# config/bb/prebuilt/etc/inittab
::sysinit:/etc/rc.d/startup

# ::respawn:/bin/sh
::respawn:/sbin/getty -n -l /bin/sh -L 115200 ttyO2 vt100

::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc.d/shutdown
