# config/bb/prebuilt/etc/inittab
# system initialisation
::sysinit:/etc/rc.d/startup

# run gettys on the serial ports
#::respawn:/sbin/getty -n -l /bin/nologin -L 57600 ttyS0 vt100
::respawn:/bin/sh

# stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc.d/shutdown
